# config/taskmaster_conf.ini
[program:ss]
command=python3 ../pgrms/shg/main_shg.py